# beacon/debian/besic-beacon.cron.d
# Start bluetooth beacon
@reboot root besic-beacon
*/5 * * * * root besic-beacon
