# autoupdates/debian/besic-autoupdates.cron.d
# BESI-C autoupdates
0 3 * * * root besic-autoupdates
