# s3upload/src/crontab
# Upload every 30 minutes
*/30 * * * * besic-s3upload
