# beacon/src/crontab
# Start bluetooth beacon
@reboot besic-beacon
